# bin/system/combine-power.service
[Unit]
 Description="transform cpu state to low-power mode"
 After=suspend.target
 
 [Service]
 ExecStart=/bin/bash 
 
 [Install]
 WantedBy=multi-user.target

